# docker/php/crontab
* * * * * /usr/local/bin/php /var/www/html/magento/bin/magento cron:run >> /var/www/html/magento/var/log/cron.log 2>&1
