# docker/php/supervisor_cron.conf
[supervisord]
nodaemon=true

[program:crond]
command=/usr/sbin/cron -f
autorestart=true
umask=022